# claude-code/rootfs/etc/services.d/claudecodeui/finish
#!/usr/bin/env bashio
# shellcheck shell=bash

declare APP_EXIT_CODE=${1}

if [[ "${APP_EXIT_CODE}" -ne 0 ]] && [[ "${APP_EXIT_CODE}" -ne 256 ]]; then
    bashio::log.warning "Claude Code UI exited with code ${APP_EXIT_CODE}"
fi

bashio::log.info "Claude Code UI service stopped, will restart..."
